# cron/fetchavgblocktime.cron
0 * * * * /usr/bin/php /var/www/sumcoinblockhalf.com/getavgblocktime.php